# roles/install_compose_application_as_service/templates/template-reload.service.j2
[Unit]
Description=Rebuild and restart {{ compose_service_name }} containers

[Service]
Type=oneshot
ExecStart=/bin/systemctl reload-or-restart {{ compose_service_name }}.service
